# corbanx-api.service
[Unit]
Description=CorbanX API - Multi-banco CLT + FGTS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/corbanx-api
ExecStart=/opt/corbanx-api/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8004 --workers 4
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
